# termux/var/service/sysinit.target/run
#!/data/data/com.termux/files/usr/bin/sh
exec systemctl3.py init 2>&1
