# var/service/systemd-journald/run
#!/data/data/com.termux/files/usr/bin/sh
exec termux-chroot systemctl3.py init 2>&1
